# nidoking/Users/manselmi/Library/LaunchAgents/com.manselmi.resticprofile.nidoking.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- vim: set ft=xml : -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/manselmi/.prefix/bin:/Users/manselmi/.prefix/sw/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>RCLONE_GID</key>
		<string>20</string><!-- staff -->
		<key>RCLONE_UID</key>
		<string>501</string><!-- manselmi -->
		<key>RESTICPROFILE_GID</key>
		<string>80</string><!-- admin -->
		<key>RESTICPROFILE_LOCK_WAIT</key>
		<string>1h</string>
		<key>RESTICPROFILE_LOG</key>
		<string>/Users/manselmi/.config/resticprofile/log/nidoking/backup.log</string>
		<key>RESTICPROFILE_NO_ANSI</key>
		<string></string>
		<key>RESTICPROFILE_NO_PRIO</key>
		<string></string>
		<key>RESTICPROFILE_UID</key>
		<string>0</string><!-- root -->
		<key>TZ</key>
		<string>UTC</string>
	</dict>
	<key>Label</key>
	<string>com.manselmi.resticprofile.nidoking.backup</string>
	<key>LowPriorityBackgroundIO</key>
	<true/>
	<key>LowPriorityIO</key>
	<true/>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-fc</string>
		<string>truncate -s 0 -- "${RESTICPROFILE_LOG}" ; exec -- fdautil exec "$(whence -p -- exec-resticprofile)" -- nidoking.backup</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>12</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
</dict>
</plist>
