# nidoking/Users/manselmi/Library/LaunchAgents/com.manselmi.resticprofile.nidoking.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- vim: set ft=xml : -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/manselmi/.prefix/bin:/Users/manselmi/.prefix/sw/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>TZ</key>
		<string>UTC</string>
	</dict>
	<key>Label</key>
	<string>com.manselmi.resticprofile.nidoking.backup</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-fc</string>
		<string>/usr/bin/truncate -s 0 -- /Users/manselmi/.config/resticprofile/log/nidoking/backup.log ; exec -- /usr/local/bin/fdautil exec /Users/manselmi/.prefix/bin/exec-resticprofile --log /Users/manselmi/.config/resticprofile/log/nidoking/backup.log -- nidoking.backup</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>12</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
</dict>
</plist>
